# dev/docker/web/config/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
directory=/usr/local/
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
redirect_stderr=true

[program:server]
command=/bin/bash -c "/usr/sbin/apache2 -D FOREGROUND"
